# packages/master/raspberry-pack/raspberry-pack.service
[Unit]
Description=Raspberry-Pack-Generic-Service
After=network-online.target

[Service]
ExecStart=/bin/bash /boot/raspberry-pack/service-starter.sh
WorkingDirectory=/boot/raspberry-pack/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target